# deploy/host/lollipop-gan.service
[Unit]
Description=Lollipop Remote GAN API
After=network.target

[Service]
Type=simple
# Suggest running as root if scripts need docker access without group tweaks, 
# or use a specific user if permissions allow.
User=root
Group=root
WorkingDirectory=/opt/lollipop-remote-gan/backend
Environment=CONFIG_FILE=/opt/lollipop-remote-gan/backend/config.yaml
ExecStart=/opt/lollipop-remote-gan/backend/lollipop-remote-gan-api
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
